# perf/uploadAndTest.timer.template
[Unit]
Description="Run upload and benchmark daily"

[Timer]
OnCalendar=*-*-* 12:00:00
Unit=uploadAndTest_PARAMS.service
Persistent=true

[Install]
WantedBy=default.target