# perf/uploadAndTest.timer.template
[Unit]
Description="Run upload and benchmark: PARAMS"

[Timer]
OnUnitActiveSec=12h
Unit=uploadAndTest_PARAMS.service
Persistent=true

[Install]
WantedBy=default.target